# systemd_unit_files/run_web_scraping.timer
[Unit]
Description=Run My Web Scraping Service every 3 minutes

[Timer]
OnCalendar=*:0/3
Persistent=true

[Install]
WantedBy=timers.target

